# src/lfs/packages/extras/spice-vdagent/spice-vdagentd
#!/bin/sh
########################################################################
# Begin spice-vdagentd
#
# Description : Spice guest agent X11 session agent daemon.
#
# Authors     : ClemaX
#
# Version     : LFS 11.3
#
# Notes       :
#
########################################################################

### BEGIN INIT INFO
# Provides:            template
# Required-Start:
# Should-Start:
# Required-Stop:
# Should-Stop:
# Default-Start:
# Default-Stop:
# Short-Description:
# Description:
# X-LFS-Provided-By:
### END INIT INFO

# shellcheck source=/dev/null
. /lib/lsb/init-functions

case "${1}" in
start)
    log_info_msg "Starting spice-vdagent daemon..."

    mkdir -p /run/spice-vdagentd

    start_daemon /usr/bin/spice-vdagentd

    evaluate_retval
    ;;

stop)
    log_info_msg "Stopping spice-vdagent daemon..."

    killproc /usr/bin/spice-vdagentd

    evaluate_retval
    ;;

restart)
    ${0} stop
    sleep 1
    ${0} start
    ;;

status)
    statusproc /usr/bin/spice-vdagentd
    ;;

*)
    echo "Usage: ${0} {start|stop|restart|status}"
    exit 1
    ;;
esac

exit 0

# End spice-vdagentd
